# deploy/valorlink-bot.service
[Unit]
Description=ValorLink Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=valorlink
Group=valorlink
WorkingDirectory=/opt/valorlink
EnvironmentFile=/opt/valorlink/.env
ExecStart=/opt/valorlink/.venv/bin/python3 bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
